# 0001_Personal/0001_Linux_Dotfiles/0000_Dell_Precision_5570/0001_NixOs/".config"/autostart/polybar.desktop
[Desktop Entry]
Exec=polybar
Icon=
Name=polybar
Path=
Terminal=True
Type=Application
